# src/program-rclive-backend.conf
[program:rclive-backend]
command = /usr/local/bin/launch-rclive-backend.sh
stdout_logfile = /var/log/rclive/forever.log
